# crontab
0 13 * * * cd "${HOME}"/repos/greymd/docker-zsh-completion && sh ./execute.sh &> /dev/null
