# crontab
* * * * * cd ${HOME}/reps/greymd/docker-zsh-completion && sh ./execute.sh &> /dev/null
